# NodeStateManager/config/nodestatemanager-daemon.service.in
#
# Copyright (C) 2012 Continental Automotive Systems, Inc.
#
# Author: Jean-Pierre.Bogler@continental-corporation.com
#
# This file is for starting the nodestatemanager-daemon as a service with systemd
#
# This Source Code Form is subject to the terms of the Mozilla Public
# License, v. 2.0. If a copy of the MPL was not distributed with this
# file, You can obtain one at http://mozilla.org/MPL/2.0/.

[Unit]
Description=NodeStateManager to provide system states, session and shutdown handling
Wants=dbus.service dlt-daemon.service pas-daemon.service
After=dbus.service dlt-daemon.service pas-daemon.service

# Start up very early, because it doesn't need anything to run prior to it:
DefaultDependencies=no

[Service]
Type=dbus
WatchdogSec=20
NotifyAccess=main
BusName=org.genivi.NodeStateManager
ExecStart=/usr/bin/NodeStateManager

[Install]
WantedBy=basic.target
